# enemy-territory/addons/usr/lib/systemd/system/etsound.service
[Unit]
Description=Enemy Territory oss sound

[Service]
ExecStart=/usr/bin/etsound

[Install]
WantedBy=multi-user.target
