# chrony/packages/chrony/postinst
#!/bin/bash

systemctl enable chronyd
systemctl enable chrony-wait
