# microcodectl/addons/etc/udev/rules.d/60-microcode.rules
# /etc/udev/rules.d/60-microcode.rules
#
# Update intel processor microcode
#

KERNEL=="microcode", RUN+="/sbin/microcode_ctl -Qu"
