# ntp/addons/etc/NetworkManager/dispatcher.d/ntpdate
#!/bin/bash

exec 1>>/var/log/ntpdate.log
exec 2>>/var/log/ntpdate.log

case "$2" in
	up)
		/usr/sbin/ntpdate -u hu.pool.ntp.org
		;;
esac
