# rsh/addons/etc/xinetd.d/exec
service exec
{
	disable     = yes
	socket_type = stream
	wait        = no
	user        = root
	server      = /usr/sbin/in.rexecd
}
