# qt5/addons/etc/profile.d/qt5.sh
QT5DIR=/usr/lib/qt5

export PATH="$PATH:$QT5DIR/bin"
export QT5DIR
