# kde-workspace/packages/kde-kdm/files
	/etc/pam.d/kde
	/etc/pam.d/kde-np
	/usr/bin/*kdm*
	/usr/lib/kde4/*kdm*.so
	/usr/lib/kde4/*kgreet*.so
	/usr/lib/kde4/libexec/kdm*
	/usr/lib/kde4/libexec/krootimage
	/usr/lib/libkworkspace.so*
	/usr/lib/systemd/system/kdm.service
	/usr/share/config/kdm*
	/usr/share/config/kdm/*
	/usr/share/doc/kdm/*
	/usr/share/kde4/services/kdm.desktop
	/usr/share/kdm/*
