# java/addons/etc/profile.d/java.sh
export PATH="$PATH:/usr/lib/java/bin:/usr/lib/java/jre/bin"
