# kernel/packages/kernel/postinst
#!/bin/sh

systemctl enable icmpredirect_fix
if [ -f /boot/grub2/grub.cfg ]; then
	mv /boot/grub2/grub.cfg /boot/grub2/grub.cfg.`date '+%F'`
fi
exec grub2-mkconfig -o /boot/grub2/grub.cfg
exec mkinitrd-installedkernel
